# ~/fixed_service/autonomous-mower.service
[Unit]
Description=Autonomous Lawn Mower Service
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/autonomous_mower
# Basic environment variables
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/autonomous_mower/src
# Run the application using the system Python with the correct module path
ExecStart=/usr/bin/python3 -m mower.main_controller
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Hardware access permissions
SupplementaryGroups=gpio i2c dialout video

[Install]
WantedBy=multi-user.target 